# package/DEBIAN/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        systemctl daemon-reload
        systemctl enable kde-collector.service
        systemctl start kde-collector.service
        ;;
esac