# systemd/kde.socket
[Unit]
Description=Kernel Docker Exporter Socket

[Socket]
ListenStream=0.0.0.0:6545
Accept=yes

[Install]
WantedBy=sockets.target